# deploy_tools/gunicorn-superlists-staging.natanielstrack.com.conf
description "Gunicorn Server for superlists-staging.natanielstrack.com"

start on net-device-up
stop on shutdown

respawn

setuid ubuntu

chdir /home/ubuntu/sites/superlists-staging.natanielstrack.com/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/superlists-staging.natanielstrack.com.socket superlists.wsgi:application
